# www/linux-c6-chrome/pkg-plist
bin/chrome
etc/cron.daily/google-chrome
opt/google/chrome/PepperFlash/libpepflashplayer.so
opt/google/chrome/PepperFlash/manifest.json
opt/google/chrome/chrome
opt/google/chrome/chrome-sandbox
opt/google/chrome/chrome_100_percent.pak
opt/google/chrome/default-app-block
opt/google/chrome/default_apps/docs.crx
opt/google/chrome/default_apps/drive.crx
opt/google/chrome/default_apps/external_extensions.json
opt/google/chrome/default_apps/gmail.crx
opt/google/chrome/default_apps/search.crx
opt/google/chrome/default_apps/youtube.crx
opt/google/chrome/google-chrome
opt/google/chrome/icudtl.dat
opt/google/chrome/libffmpegsumo.so
opt/google/chrome/libwidevinecdm.so
opt/google/chrome/libwidevinecdmadapter.so
opt/google/chrome/locales/am.pak
opt/google/chrome/locales/ar.pak
opt/google/chrome/locales/bg.pak
opt/google/chrome/locales/bn.pak
opt/google/chrome/locales/ca.pak
opt/google/chrome/locales/cs.pak
opt/google/chrome/locales/da.pak
opt/google/chrome/locales/de.pak
opt/google/chrome/locales/el.pak
opt/google/chrome/locales/en-GB.pak
opt/google/chrome/locales/en-US.pak
opt/google/chrome/locales/es-419.pak
opt/google/chrome/locales/es.pak
opt/google/chrome/locales/et.pak
opt/google/chrome/locales/fa.pak
opt/google/chrome/locales/fi.pak
opt/google/chrome/locales/fil.pak
opt/google/chrome/locales/fr.pak
opt/google/chrome/locales/gu.pak
opt/google/chrome/locales/he.pak
opt/google/chrome/locales/hi.pak
opt/google/chrome/locales/hr.pak
opt/google/chrome/locales/hu.pak
opt/google/chrome/locales/id.pak
opt/google/chrome/locales/it.pak
opt/google/chrome/locales/ja.pak
opt/google/chrome/locales/kn.pak
opt/google/chrome/locales/ko.pak
opt/google/chrome/locales/lt.pak
opt/google/chrome/locales/lv.pak
opt/google/chrome/locales/ml.pak
opt/google/chrome/locales/mr.pak
opt/google/chrome/locales/ms.pak
opt/google/chrome/locales/nb.pak
opt/google/chrome/locales/nl.pak
opt/google/chrome/locales/pl.pak
opt/google/chrome/locales/pt-BR.pak
opt/google/chrome/locales/pt-PT.pak
opt/google/chrome/locales/ro.pak
opt/google/chrome/locales/ru.pak
opt/google/chrome/locales/sk.pak
opt/google/chrome/locales/sl.pak
opt/google/chrome/locales/sr.pak
opt/google/chrome/locales/sv.pak
opt/google/chrome/locales/sw.pak
opt/google/chrome/locales/ta.pak
opt/google/chrome/locales/te.pak
opt/google/chrome/locales/th.pak
opt/google/chrome/locales/tr.pak
opt/google/chrome/locales/uk.pak
opt/google/chrome/locales/vi.pak
opt/google/chrome/locales/zh-CN.pak
opt/google/chrome/locales/zh-TW.pak
opt/google/chrome/nacl_helper
opt/google/chrome/nacl_helper_bootstrap
opt/google/chrome/nacl_irt_x86_32.nexe
opt/google/chrome/natives_blob.bin
opt/google/chrome/product_logo_128.png
opt/google/chrome/product_logo_16.png
opt/google/chrome/product_logo_22.png
opt/google/chrome/product_logo_24.png
opt/google/chrome/product_logo_256.png
opt/google/chrome/product_logo_32.png
opt/google/chrome/product_logo_32.xpm
opt/google/chrome/product_logo_48.png
opt/google/chrome/product_logo_64.png
opt/google/chrome/resources.pak
opt/google/chrome/snapshot_blob.bin
opt/google/chrome/xdg-mime
opt/google/chrome/xdg-settings
usr/bin/google-chrome-stable
usr/share/applications/google-chrome.desktop
usr/share/gnome-control-center/default-apps/google-chrome.xml
usr/share/man/man1/google-chrome.1
