# lpjs/files/lpjs_dispatchd.in
#!/bin/sh

# PROVIDE: lpjs_dispatchd
# REQUIRE: DAEMON munge
# BEFORE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# lpjs_dispatchd_enable (bool):   Set to NO by default.
#               Set it to YES to enable lpjs_dispatchd.
#

. /etc/rc.subr

name="lpjs_dispatchd"
rcvar=lpjs_dispatchd_enable

pidfile=/var/run/$name.pid

load_rc_config $name

: ${lpjs_dispatchd_enable="NO"}

start_cmd=lpjs_dispatchd_start
stop_cmd=lpjs_dispatchd_stop

lpjs_dispatchd_start() {
    checkyesno lpjs_dispatchd_enable && echo "Starting $name." && \
	%%PREFIX%%/sbin/$name --daemonize $lpjs_dispatchd_flags
}

lpjs_dispatchd_stop() {
    if [ -e $pidfile ]; then
	checkyesno lpjs_dispatchd_enable && echo "Stopping $name." && \
	    kill `cat $pidfile`
    else
	echo "Stopping $name." && pkill $name
    fi
}

run_rc_command "$1"
