# net/wifi_select/files/wifi_select.in
#!/bin/sh
#
# $FreeBSD$
#

# FIXME:
# Adding BEFORE: ntpdate would be nice so that ntpdate will actually
# work on a wireless laptop, but this causes a circular dependency in
# rcorder

# PROVIDE: wifi_select
# REQUIRE: DAEMON
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable `wifi_select':
#
# wifi_select_enable="YES"
# wifi_select_flags="<set as needed>"
#
# See wifi_select(1) for wifi_select_flags
#

. /etc/rc.subr

name="wifi_select"
rcvar=`set_rcvar`

pidfile="/var/run/$name.pid"
command="%%PREFIX%%/libexec/wifi_select_daemon"
command_args=""
# required_files="%%PREFIX%%/etc/wifi_select/stations.conf"

start_cmd=wifi_select_start
stop_cmd=wifi_select_stop

wifi_select_daemon_bin=%%PREFIX%%/libexec/wifi_select_daemon

wifi_select_start() {
    checkyesno wifi_select_enable && echo "Starting wifi_select." && \
	${wifi_select_daemon_bin}
}

wifi_select_stop()
{
    checkyesno wifi_select_enable && echo "Stopping wifi_select." && \
	${wifi_select_daemon_bin} -k
}

# read configuration and set defaults
load_rc_config "$name"
# : ${wifi_select_enable="NO"}

run_rc_command "$1"

